# systemd.socket
[Socket]
ListenStream=/run/api-linx.sk

[Install]
WantedBy=sockets.target
